# config/client/services.conf
# Required configuration files structure:

# 1. Client configuration files (place in config/client/)
# ------------------------------------------------------

## config/client/services.conf
[program:init]
command=/app/init.sh
autostart=true
autorestart=false
priority=5
startretries=1
startsecs=0
user=root
stdout_logfile=/var/log/init.log
stderr_logfile=/var/log/init_err.log

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 %(ENV_SCREEN_RESOLUTION)s
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb_err.log

[program:x11vnc]
command=/usr/bin/x11vnc -xkb -noxrecord -noxfixes -noxdamage -display :1 -forever -shared -passwd %(ENV_VNC_PASSWORD)s
autostart=true
autorestart=true
priority=20
depends_on=xvfb
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc_err.log

[program:novnc]
command=sh -c 'sleep 3 && /usr/bin/python3 /app/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 6080 --web /app/noVNC'
autostart=true
autorestart=true
priority=30
depends_on=x11vnc
startsecs=10
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc_err.log

[program:09launcher]
command=sh -c 'sleep 5 && /usr/bin/java -jar /app/09launcher.jar --debug'
directory=/home/%(ENV_USER_NAME)s
autostart=true
autorestart=true
priority=40
user=%(ENV_USER_NAME)s
depends_on=xvfb
stdout_logfile=/var/log/09launcher.log
stderr_logfile=/var/log/09launcher_err.log
startretries=5
stopasgroup=true
killasgroup=true
startsecs=10
environment_file=/etc/supervisor/env/09launcher.env

[group:vnc-client]
programs=init,xvfb,x11vnc,novnc,09launcher
priority=999

## config/client/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
user=root
childlogdir=/var/log/supervisor

[include]
files = /etc/supervisor/conf.d/*.conf

# 2. Server configuration files (place in config/server/)
# ------------------------------------------------------

## config/supervisord.conf (for server)
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
user=root
childlogdir=/var/log/supervisor

[include]
files = /etc/supervisor/conf.d/*.conf

; MySQL database service with proper command
[program:mysql]
command=/usr/bin/mysqld_safe --datadir=/var/lib/mysql
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/mysql.log
stderr_logfile=/var/log/mysql_err.log
startretries=5
startsecs=10
stopwaitsecs=120
stdout_events_enabled=true
stderr_events_enabled=true

; Database initialization - now with mysql health check
[program:db-init]
command=/app/scripts/db-init.sh
autostart=true
autorestart=false
startsecs=0
startretries=5
priority=20
environment=MYSQL_USER="%(ENV_MYSQL_USER)s",MYSQL_PASSWORD="%(ENV_MYSQL_PASSWORD)s",MYSQL_DATABASE="%(ENV_MYSQL_DATABASE)s",MYSQL_HOST="localhost"
stdout_logfile=/var/log/db-init.log
stderr_logfile=/var/log/db-init_err.log
exitcodes=0

; 2009scape game server with dependency on db-init
[program:2009scape-server]
command=/app/scripts/wait-for-db-init.sh /app/scripts/start-server.sh
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/2009scape-server.log
stderr_logfile=/var/log/2009scape-server_err.log
environment=MYSQL_HOST="%(ENV_MYSQL_HOST)s",MYSQL_USER="%(ENV_MYSQL_USER)s",MYSQL_PASSWORD="%(ENV_MYSQL_PASSWORD)s",MYSQL_DATABASE="%(ENV_MYSQL_DATABASE)s",JAVA_OPTS="%(ENV_JAVA_OPTS)s"
startretries=5
stopasgroup=true
killasgroup=true
startsecs=10

; 2009scape management interface with dependency check
[program:2009scape-management]
command=sh -c "until supervisorctl status 2009scape-server | grep -q 'RUNNING'; do echo 'Waiting for 2009scape server to start...'; sleep 5; done && cd /app/2009scape && ./management/management.sh"
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/2009scape-management.log
stderr_logfile=/var/log/2009scape-management_err.log
environment=MYSQL_HOST="%(ENV_MYSQL_HOST)s",MYSQL_USER="%(ENV_MYSQL_USER)s",MYSQL_PASSWORD="%(ENV_MYSQL_PASSWORD)s",MYSQL_DATABASE="%(ENV_MYSQL_DATABASE)s",JAVA_OPTS="%(ENV_JAVA_OPTS)s"
startretries=3
stopasgroup=true
killasgroup=true
startsecs=10

; Group all 2009scape services
[group:2009scape]
programs=mysql,db-init,2009scape-server,2009scape-management
priority=999

# 3. MySQL initialization (place in config/mysql/init/)
# ----------------------------------------------------

## config/mysql/init/01-initialize-database.sql
-- Create user and grant privileges
CREATE USER IF NOT EXISTS '${MYSQL_USER}'@'%' IDENTIFIED BY '${MYSQL_PASSWORD}';
GRANT ALL PRIVILEGES ON `${MYSQL_DATABASE}`.* TO '${MYSQL_USER}'@'%';
FLUSH PRIVILEGES;

-- Create basic tables for 2009scape
CREATE TABLE IF NOT EXISTS `players` (
  `id` int(11) NOT NULL AUTO_INCREMENT,
  `username` varchar(255) NOT NULL,
  `password` varchar(255) NOT NULL,
  `credits` int(11) DEFAULT 0,
  `client_rights` int(11) DEFAULT 0,
  `donator_credits` int(11) DEFAULT 0,
  `friends` text DEFAULT NULL,
  `ignores` text DEFAULT NULL,
  PRIMARY KEY (`id`),
  UNIQUE KEY `username` (`username`)
) ENGINE=InnoDB DEFAULT CHARSET=utf8mb4;

CREATE TABLE IF NOT EXISTS `game_logs` (
  `id` int(11) NOT NULL AUTO_INCREMENT,
  `player_id` int(11) NOT NULL,
  `timestamp` timestamp NOT NULL DEFAULT current_timestamp(),
  `action` varchar(255) NOT NULL,
  `details` text DEFAULT NULL,
  PRIMARY KEY (`id`),
  KEY `player_id` (`player_id`),
  CONSTRAINT `game_logs_ibfk_1` FOREIGN KEY (`player_id`) REFERENCES `players` (`id`)
) ENGINE=InnoDB DEFAULT CHARSET=utf8mb4;